# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
 
[program:ssh]
command=/root/ssh.sh

[program:vnc]
command=/root/vnc.sh